# numlock-off.service
[Unit]
Description=Force NumLock off — polling loop
After=graphical-session.target

[Service]
Type=simple
ExecStart=/bin/bash -c 'while true; do DISPLAY=:0 /usr/bin/numlockx off; sleep 30; done'
ExecStop=/usr/bin/pkill -f "numlockx"
Restart=always
RestartSec=10
StartLimitIntervalSec=60
StartLimitBurst=3

[Install]
WantedBy=default.target
